# recipes-connectivity/bluetooth/files/broadcom-bluetooth/brcm-btfw-load@.service
[Unit]
Description=Load Bluetooth firmware

[Service]
Type=oneshot
ExecStartPre=/usr/sbin/rfkill unblock bluetooth
ExecStart=/bin/bash /usr/bin/brcm-btfw-load.sh %I
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
